# infra/systemd/work-planner-production.service
[Unit]
Description=Work Planner API (production)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=work-planner
Group=work-planner
WorkingDirectory=/var/www/work-planner/production/backend/current
EnvironmentFile=/etc/work-planner/production.env
ExecStart=/usr/bin/node apps/api/dist/server.js
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/var/www/work-planner/production/data /var/www/work-planner/production/backups /var/www/work-planner/production/backend/current

[Install]
WantedBy=multi-user.target

